# deploy/systemd/uma-sync.service
[Unit]
Description=poly_uma uma-sync service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=uma-sync
Group=uma-sync
WorkingDirectory=/opt/uma-sync
EnvironmentFile=/opt/uma-sync/.uma-sync.env
Environment=API_ADDR=0.0.0.0:8011
ExecStart=/opt/uma-sync/uma-sync
Restart=always
RestartSec=3

# 基础安全加固（按需放开）
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target

